# anki_screenshot_creator/launchd/com.anki-screenshot-creator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- TEMPLATE: install.sh fills in PYTHON_PATH and PROJECT_DIR before loading -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.anki-screenshot-creator</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__PROJECT__/flask_server.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/anki-screenshot-creator.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/anki-screenshot-creator.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
